# deploy/dgx-spark/agentics-web.service
[Unit]
Description=Agentics web frontend
After=network-online.target agentics-api.service
Wants=network-online.target

[Service]
Type=simple
User=agentics
Group=agentics
WorkingDirectory=/opt/agentics/current/frontends/web
EnvironmentFile=/etc/agentics/agentics.env
Environment=PATH=/opt/agentics/current/bin:/opt/agentics/current/frontends/web/node_modules/.bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/opt/agentics/current/bin/bun run start -- -p ${AGENTICS_WEB_PORT}
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/agentics/current/frontends/web/.next/cache

[Install]
WantedBy=multi-user.target
